# webapp/deploy/macroquest-webapp.service
[Unit]
Description=MacroQuest Admin Panel
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=mqadmin
Group=mqadmin
WorkingDirectory=/opt/mq-admin
EnvironmentFile=/opt/mq-admin/.env
ExecStart=/opt/mq-admin/venv/bin/uvicorn webapp.app:app --host 127.0.0.1 --port 8000 --workers 1
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mq-admin
NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=/opt/mq-admin/webapp/data /opt/mq-admin/.env

[Install]
WantedBy=multi-user.target
