# webapp/deploy/macroquest-webapp.service
[Unit]
Description=MacroQuest Admin Panel (Build Pipeline)
Documentation=https://github.com/teichoui/macroquest
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User=mqadmin
Group=mqadmin
WorkingDirectory=/opt/mq-admin

# Load environment variables from .env file
EnvironmentFile=/opt/mq-admin/.env

# Uvicorn application entry point
ExecStart=/opt/mq-admin/venv/bin/uvicorn webapp.app:app \
    --host 127.0.0.1 \
    --port 8000 \
    --workers 2 \
    --no-access-log \
    --log-level warning

# Graceful shutdown: send SIGTERM and wait up to 10s before SIGKILL
ExecStop=/bin/kill -TERM $MAINPID
TimeoutStopSec=10

# Restart policy
Restart=always
RestartSec=3
StartLimitIntervalSec=60
StartLimitBurst=5

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/mq-admin/webapp/data
ReadWritePaths=/opt/mq-admin/webapp/logs

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mq-admin

[Install]
WantedBy=multi-user.target
